# salt/slack-irc/config/slack-irc.upstart.conf.jinja
description "slack-irc bridge"
start on runlevel [2345]
stop on runlevel [06]

console log
setuid slack-irc
setgid slack-irc

respawn

exec slack-irc --config /etc/slack-irc.json
